# guard_tunnel_bot.service
[Unit]
Description=Guard Tunnel Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/guard-tunnel-bot
EnvironmentFile=/opt/guard-tunnel-bot/.env
ExecStart=/opt/guard-tunnel-bot/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target